# OS/alpine/argononed.start
#!/sbin/openrc-run
command=/usr/sbin/argononed
command_args=""
pidfile=/run/argononed.pid
description="Argon One Fan and Button Daemon Service"

depend() {
        need localmount
        after modules
        after rpi3-config-mv
}
